# supervisor_conf/supervisord.conf
[program:BootAssist]
command=python3 /home/pi/ColdBootAssistant/main.py
autorestart=true